# scripts/reactapp.conf
#!upstart
description "node.js server"
author      "ssridharan"

env NAME=reactapp

start on stopped rc RUNLEVEL=[2345]
stop on runlevel [!2345]

# Automatically Respawn:
respawn
respawn limit 99 5

script
    export HOME="/root"
    export USER="reactappuser"

    sudo -u $USER APP_ENV=$ENVIRONEMT /usr/bin/node /apps/$NAME/server/server.js >> /logs/$NAME/$NAME.log 2>&1

end script

pre-start script
    # Date format same as (new Date()).toISOString() for consistency
    echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Starting" >> /logs/$NAME/$NAME.log
end script

pre-stop script
    rm /var/run/$NAME.pid
    echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Stopping" >> /var/log/$NAME.log
end script
